# vagrant/supervisord.conf
[program:grader]
command = /usr/bin/php /var/grader/client.php grader
process_name = %(program_name)s_%(process_num)02d
numprocs = 5
exitcodes = 0
stopasgroup = true
directory = /var/grader/